# src/build/modules/logs/files/etc/supervisor/conf.d/dockergen.conf
[program:dockergen]
command = /usr/local/bin/docker-gen -watch -only-exposed -notify "/usr/bin/supervisorctl restart fluentd" /root/fluentd.conf.tmpl /etc/td-agent/td-agent.conf
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
