# src/build/modules/logs/files/etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command = /opt/elasticsearch/bin/elasticsearch
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
