# src/build/modules/logs/files/etc/supervisor/conf.d/fluentd.conf
[program:fluentd]
command = /usr/sbin/td-agent
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
