# src/build/modules/logs/files/etc/supervisor/conf.d/kibana.conf
[program:kibana]
command = /opt/kibana/bin/kibana
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
